# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:apache]
# Waits for MariaDB on 127.0.0.1:5054, then runs the official WordPress
# entrypoint (generates wp-config.php from env) + apache2 in the foreground.
command=/usr/local/bin/start-apache.sh
priority=20
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
